# scripts/service/transmission/run
#!/usr/bin/with-contenv sh

sleep 15
 
exec s6-setuidgid transmission /usr/bin/transmission-daemon --foreground --config-dir /etc/transmission-daemon