# scripts/service/openvpn/run
#!/usr/bin/with-contenv sh

openvpn --cd /etc/openvpn --config service.conf --auth-nocache