# scripts/service/openvpn/run
#!/usr/bin/with-contenv sh

s6-setuidgid openvpn openvpn --cd /etc/openvpn --config service.conf --auth-nocache